# .config/systemd/user/wall-agenda.service
[Unit]
Description=Render wallpaper with next calendar event

[Service]
Type=oneshot
ExecStart=%h/.dotfiles/scripts/wall_agenda.sh
