# systemd/user/vigil.service
[Unit]
Description=Vigil Wayland desktop cockpit
After=graphical-session.target
PartOf=graphical-session.target
Conflicts=dunst.service

[Service]
Type=simple
ExecStart=%h/.local/bin/vigil
Restart=on-failure
RestartSec=2
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
Environment=XDG_CONFIG_HOME=%h/.config

[Install]
WantedBy=graphical-session.target
